# pn-config/papyri.info
#!/bin/bash
# chkconfig: 345 86 16
# description: Apache Tomcat servlet container
#------------------------------------------------------------------------------
# $Id$
#
# papyri.info control script
#
#------------------------------------------------------------------------------
export JAVA_HOME=/usr/lib/jvm/jre-1.8.0-openjdk.x86_64
export JAVA_TOOL_OPTIONS="-Dfile.encoding=UTF8"
CATALINA_SH=/usr/local/tomcat-sosol/bin/catalina.sh
PN_SH=/usr/local/tomcat-pn/bin/catalina.sh
export PN_PID=/usr/local/tomcat-pn/bin/tomcat.pid

starttc () {
  mv /usr/local/tomcat-sosol/logs/catalina.out /usr/local/tomcat-sosol/logs/catalina.$(date +"%Y-%m-%d.%H.%M.%S").log
  export CATALINA_PID=/usr/local/tomcat-sosol/bin/tomcat.pid
  export JAVA_OPTS="-server -Xms1500m -Xmx4G -XX:MaxPermSize=256m -XX:+HeapDumpOnOutOfMemoryError -XX:HeapDumpPath=/srv/data/papyri.info/heap-dumps -XX:+UseConcMarkSweepGC -XX:+CMSIncrementalMode -XX:+CMSClassUnloadingEnabled -verbose:gc -Xloggc:gc.log -XX:+PrintGCDetails -XX:+PrintGCTimeStamps -XX:+PrintGCDateStamps -XX:+PrintTenuringDistribution -XX:+PrintHeapAtGC -Djruby.objectspace.enabled=false -Djruby.thread.pool.enabled=true -Djruby.compat.version=2.0"
  su -s /bin/sh tomcat -c "$CATALINA_SH start"
}

stoptc () {
  export CATALINA_PID=/usr/local/tomcat-sosol/bin/tomcat.pid
  su -s /bin/sh tomcat -c "$CATALINA_SH stop"
  sleep 10
  if [ -e "$CATALINA_PID" ]
  then
          kill -TERM `cat $CATALINA_PID`
          sleep 10
          if [ -e "$CATALINA_PID" ]
          then
                  kill -9 `cat $CATALINA_PID`
          fi
  fi
}

startpn() {
  cd /usr/local/tomcat-pn/logs
  export JAVA_OPTS="-server -Xms500m -Xmx1G"
  export CATALINA_PID=/usr/local/tomcat-pn/bin/tomcat.pid
  su -s /bin/sh tomcat -c "$PN_SH start"
}

stoppn () {
  export CATALINA_PID=/usr/local/tomcat-pn/bin/tomcat.pid
  su -s /bin/sh tomcat -c "$PN_SH stop"
  sleep 10
  if [ -e "$PN_PID" ]
  then
          kill -TERM `cat $PN_PID`
          sleep 10
          if [ -e "PN_PID" ]
          then
                  kill -9 `cat $PN_PID`
          fi
  fi
}

startsolr() {
  export PATH=$PATH:/usr/sbin
  su -s /bin/sh tomcat -c "/usr/local/solr/bin/solr start -m 2G -s /srv/data/papyri.info/solr"
}

stopsolr() {
  su -s /bin/sh tomcat -c "/usr/local/solr/bin/solr stop -all"
}

startjena () {
  cd /usr/local/jena/jena-fuseki
  su -s /bin/sh tomcat -c "./fuseki-server --update --port 8090 --loc=/srv/data/papyri.info/jena /pi &> /usr/local/jena/jena.log" &
}

stopjena () {
  kill -15 `pgrep -nf fuseki`
}

case "$1" in
'start')
  startjena
  echo "waiting 5 seconds for fuseki to start..."
  sleep 6
  startsolr
  startpn
  starttc
  ;;

'stop')
  stoppn
  stoptc
  stopsolr
  stopjena
  ;;

'restart')
  stoppn
  stoptc
  stopsolr
  stopjena
  sleep 20
  startjena
  echo "waiting 5 seconds for fuseki to start..."
  sleep 6
  startsolr
  startpn
  starttc
  ;;

'start-tc')
  starttc
  ;;

'stop-tc')
  stoptc
  ;;

'start-pn')
  startpn
  ;;

'stop-pn')
  stoppn
  ;;

'start-solr')
  startsolr
  ;;

'stop-solr')
  stopsolr
  ;;

'restart-solr')
  stopsolr
  sleep 10
  startsolr
  ;;

'start-jena')
  startjena
  ;;

'stop-jena')
  stopjena
  ;;

*)
  echo "Usage: $0 { start | stop | restart | start-tc | stop-tc | start-pn | stop-pn | start-solr | stop-solr | restart-solr | start-jena | stop-jena }"
  exit 1
  ;;
esac


#------------------------------------------------------------------------------
# END OF SCRIPT
#------------------------------------------------------------------------------